# systemd/mnemon.timer
# Fires mnemon.service every 5 minutes. Install alongside the service and
# enable with: sudo systemctl enable --now mnemon.timer
[Unit]
Description=Run MNEMON ingestion every 5 minutes

[Timer]
# Every 5 min on the wall clock. Minute-based, so the host timezone is
# irrelevant — MNEMON buckets timestamps in UTC internally. The run itself
# decides which jobs are due, so most ticks only fetch market_state.
OnCalendar=*:0/5
# Run a missed cycle once on boot if the VPS was down at the scheduled time.
Persistent=true
# Let systemd spread load slightly rather than firing dead-on the second.
AccuracySec=30s

[Install]
WantedBy=timers.target
